# etc/supervisord.ini
[program:kursobot]
command=/usr/local/kursobot/bin/app -config=/mnt/kursobot/config/kursobot.toml
autostart=true
autorestart=true
startretries=5
user=kursobot
redirect_stderr=true
stdout_logfile=/mnt/kursobot/log/stdout.log
stderr_logfile=/mnt/kursobot/log/stderr.log

[supervisord]
pidfile=/tmp/supervisord.pid
nodaemon=false
logfile=/mnt/kursobot/log/supervisord.log
loglevel=info
logfile_maxbytes=50MB
logfile_backups=10
childlogdir=/mnt/kursobot/log/supervisord/
minfds=2048
minprocs=256